# cmd/masky-server/masky-server.service
[Unit]
Description=Masky Server
After=network.target

[Service]
Type=exec
ExecStart=/usr/local/bin/masky-server
